# current/slocate/SOURCES/slocate.cron
#!/bin/sh

nice --adjustment=+19 /usr/bin/slocate -c -u 
