# removed/releases/1.1-CURRENT/clockspeed/SOURCES/clockspeed.cron
#!/bin/sh

[ -f /etc/sysconfig/clockspeed ] && . /etc/sysconfig/clockspeed

if [ "$NTPSERVER" != "" ]; then
    sntpclock $NTPSERVER > /var/lib/clockspeed/adjust 2>/dev/null
fi

